# src/usr/local/etc/rc.d/frp_ssserver
#!/bin/sh

# PROVIDE: frp_ssserver
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="frp_ssserver"
rcvar="frp_ssserver_enable"

load_rc_config $name

: ${frp_ssserver_enable:="NO"}

pidfile="/var/run/frp_ssserver.pid"
logdir="/var/log/frp"
confdir="/usr/local/etc/frp"
config="${confdir}/ss.json"
binary="/usr/local/bin/ssserver"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
restart_cmd="${name}_restart"

frp_ssserver_start()
{
    if [ ! -d "$logdir" ]; then
        mkdir -p "$logdir"
    fi

    # Generate configs via template engine
    /usr/local/bin/pluginctl -c frp

    if [ ! -x "$binary" ]; then
        echo "ssserver binary not found: $binary"
        return 1
    fi

    if [ ! -f "$config" ]; then
        echo "Shadowsocks config not found: $config"
        return 1
    fi

    echo "Starting Shadowsocks server..."
    /usr/sbin/daemon -p "$pidfile" -f "$binary" -c "$config"
}

frp_ssserver_stop()
{
    if [ -f "$pidfile" ]; then
        pid=$(cat "$pidfile")
        if kill -0 "$pid" 2>/dev/null; then
            echo "Stopping Shadowsocks server..."
            kill "$pid"
            count=0
            while kill -0 "$pid" 2>/dev/null && [ $count -lt 10 ]; do
                sleep 0.5
                count=$((count + 1))
            done
            if kill -0 "$pid" 2>/dev/null; then
                kill -9 "$pid"
            fi
        fi
        rm -f "$pidfile"
    fi
}

frp_ssserver_status()
{
    if [ -f "$pidfile" ]; then
        pid=$(cat "$pidfile")
        if kill -0 "$pid" 2>/dev/null; then
            echo "Shadowsocks server is running (PID: ${pid})"
            return 0
        fi
    fi
    echo "Shadowsocks server is not running"
    return 1
}

frp_ssserver_restart()
{
    frp_ssserver_stop
    sleep 1
    frp_ssserver_start
}

run_rc_command "$1"
